# devel/.private/rosconsole/etc/catkin/profile.d/10.rosconsole.sh
# generated from rosconsole/env-hooks/10.rosconsole.sh.develspace.em

export ROSCONSOLE_CONFIG_FILE="/home/madhav/Documents/ros1_melodic_ws/src/rosconsole/config/rosconsole.config"
