# devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

. "/home/madhav/Documents/ros1_melodic_ws/src/ros/tools/rosbash/roszsh"
