# devel/etc/catkin/profile.d/15.rosbash.fish
/home/madhav/Documents/ros1_melodic_ws/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.fish